# 02-Areas/Infrastructure/pi-cross-node-comms/ansible/systemd/fleet-cooling.service
[Unit]
Description=Apply fleet cooling settings (fans + CPU governor)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/friasc/pi-cross-node-comms/fleet-cooling-setup.sh
User=root

[Install]
WantedBy=multi-user.target
