# 03-Resources/Infrastructure/technical-infrastructure-legacy/ansible/templates/decompose-watcher.service.j2
[Unit]
Description=Decomposition Watcher — Auto-process complexity triggers
After=network.target

[Service]
Type=simple
User={{ ansible_user | default('friasc') }}
WorkingDirectory={{ playbook_dir | dirname | dirname | dirname }}
ExecStart=python3 technical-infrastructure/scripts/decompose-watcher.py --daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
