# 03-Resources/Infrastructure/technical-infrastructure-legacy/ansible/templates/decompose-watcher.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingdesk.decompose-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>python3</string>
        <string>{{ repo_root }}/technical-infrastructure/scripts/decompose-watcher.py</string>
        <string>--daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ ansible_env.HOME }}/.pi/decomposition-triggers/watcher.log</string>
    <key>StandardErrorPath</key>
    <string>{{ ansible_env.HOME }}/.pi/decomposition-triggers/watcher.error.log</string>
</dict>
</plist>
