# 03-Resources/Infrastructure/technical-infrastructure-legacy/ansible/templates/gist-worker@.service.j2
[Unit]
Description=Gist Message Queue Worker for node %i
After=network.target

[Service]
Type=oneshot
User=root
Environment="GIST_WORKER_GIST_ID={{ gist_worker_gist_id }}"
Environment="GIST_WORKER_NODE=%i"
Environment="GIST_WORKER_INTERVAL={{ gist_worker_poll_interval }}"
{% if gist_worker_github_token is defined %}
Environment="GITHUB_TOKEN={{ gist_worker_github_token }}"
{% endif %}
ExecStart=/usr/local/bin/gist-worker.py --check-once
StandardOutput=journal
StandardError=journal
